# scripts/launchd/com.agent-zero.git-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.agent-zero.git-sync</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/deepak/hackathon/agent-zero/scripts/automation/git_sync.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/deepak/hackathon/agent-zero</string>

  <key>StartInterval</key>
  <integer>300</integer>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/deepak/.nvm/versions/node/v20.19.5/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/deepak/hackathon/agent-zero/.codex-automation/logs/git-sync.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/deepak/hackathon/agent-zero/.codex-automation/logs/git-sync.err.log</string>
</dict>
</plist>
